# bootstrap/roles/ldap/files/sssd_restart7.cron
#restart sssd#############################################
# sssd seems to fail periodically. setup cron to restart
0 0 * * * root /bin/systemctl stop sssd.service >/dev/null 2>&1 ; rm -f /var/lib/sss/db/*.ldb >/dev/null 2>&1 ; /bin/systemctl start sssd.service >/dev/null 2>&1
