# bootstrap/roles/ldap/files/sssd_restart.cron
#restart sssd#############################################
# sssd seems to fail periodically. setup cron to restart
0 0 * * * root /etc/init.d/sssd stop >/dev/null 2>&1 ; rm -f /var/lib/sss/db/*.ldb >/dev/null 2>&1 ; /etc/init.d/sssd start >/dev/null 2>&1
